# build/supervisord.conf
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL for a unix socket
username = dummy
password = dummy

[supervisord]
nodaemon=true
logfile=/dev/stderr ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; (AUTO child log dir, default $TEMP)
logfile_maxbytes=0
; It resolves the «CRIT Supervisor running as root (no user in config file)» warning in the log.
user = root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[unix_http_server]
file=/var/run/supervisor.sock ; (the path to the socket file)
chmod=0700 ; sockef file mode (default 0700)
username = dummy
password = dummy


[program:apache]
command=/bin/bash -c 'pidof apache2 apache2ct|wc -c|grep ^0 -q && killall $(pidof apache2 apache2ctl 2>/dev/null);grep apache /proc/$(cat /var/run/apache2/apache2.pid)/cmdline -q || rm /var/run/apache2/apache2.pi d; apache2ctl -DFOREGROUND;sleep 0.2'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true

[program:cron]
command=cron -f
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[include]
files = /etc/supervisor/conf.d/*.conf


